# services/sonarr.service.conf
[program:sonarr]
command=mono /opt/sonarr/NzbDrone.exe
user=docker
stopwaitsecs=5
stdout_events_enabled = true
stderr_events_enabled = true
